# c.market-data/crontab
* */1 * * * sh /jobs/put_publish_latest_price_information.sh
