# c.exchanges/crontab
*/5 * * * * sh /jobs/get_binance_exchange_status.sh
*/5 * * * * sh /jobs/get_bitfinex_exchange_status.sh
